# upstart/vrm-retention-online.conf
# upstart/vrm-retention-online.conf
#
description "vrm-retention-online"

start on runlevel [345]
stop on starting shutdown

respawn

script
    cd /opt/vrm-retention-online
    echo $$ > /var/run/vrm-retention-online.pid

    exec su -s /bin/bash -c 'exec /opt/vrm-retention-online/bin/vrm-retention-online  -Dconfig.file=/opt/vrm-retention-online/conf/application.conf -Dlogger.file=/opt/vrm-retention-online/conf/prod-logger.xml' vrm-retention-online
end script

pre-start script
    rm -f /opt/vrm-retention-online/RUNNING_PID
end script
